# templates/etc/systemd/system/numad.service.d/puppet.conf.epp
<%- |
  Stdlib::Absolutepath $numad_binary = '/usr/sbin/numad',
  Array $numad_arguments = []
| -%>
[Service]
ExecStart=
ExecStart=<%= $numad_binary %> <%= $numad_arguments.join(' ') %>
